# etc/cron.d/wpwatch
* * * * *	root	/usr/local/bin/wpwatch.cron

