# usr/local/bin/wpwatch.cron
#!/usr/bin/bash
# This script checks inotifywait to watch for changes wordpress files
DATE=`/bin/date`
HOSTNAME=`hostname`
if [[ `wc -w /var/log/wpwatch.log | cut -f 1 -d ' '` > 0 ]] ; then 
  echo 'blah'
  /usr/bin/cat /var/log/wpwatch.log | /usr/bin/mailx admin@netzpolitik.org -a 'FROM:mwa@netzpolitik.org' -a "SUBJECT:!WPALARM!$HOSTNAME"
  # reset log file
  /usr/bin/mv /var/log/wpwatch.log /var/log/wpwatch.log$DATE
  touch  /var/log/wpwatch.log
fi
exit 0
